# graphios.conf
description "Graphios"
author  "author"

start on runlevel [234]
stop on runlevel [016]

exec /opt/nagios/bin/graphios.py
