# init/ubuntu/graphios.conf
description "Graphios"
author  "Shawn Sterling"

start on runlevel [234]
stop on runlevel [016]

exec /usr/local/bin/graphios.py
